# git2consul.timer
[Timer]
OnCalendar=*:0/2
Unit=git2consul.service
